# vnc.service
[Unit]
Description=vnc script

[Service]
Type=simple
WorkingDirectory=/
ExecStart=bash /entrypoint.sh

[Install]
WantedBy=multi-user.target